# services/playground/run
#! /bin/bash
echo 'Starting playground'
/playgroundgo/src/github.com/Serulian/playground/playground
echo 'Playground exited'
